# backups01/etc/cron.d/containers
# Container backups once daily at 2am
0 2 * * *  root    /root/container_download.py
