# backups01/etc/cron.d/dbbackup
MAILTO=george.hicks@hpe.com

10 0,6,18 * * *	root	/data/scripts/dbbackup.sh
0 23 * * *	root	/data/scripts/rotate_sql.sh

