# dad/etc/cron.d/update-hpcloud-git-mirrors
#0,30 * * * *	root	/usr/local/bin/update-hpcloud-git-mirrors
0,30 * * * *	www-data	/usr/local/bin/update-hpcloud-git-mirrors

