# dad/etc/cron.d/update-activestate-git-mirrors
*/5 * * * *	www-data	/usr/local/bin/update-activestate-git-mirrors

